# easy/easy.service
[Unit]
Description=easy time probe
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/easy
Environment=EASY_DB=/var/lib/easy/easy.db
Environment=EASY_PORT=8080
Environment=EASY_INTERVAL=10s
ExecStart=/opt/easy/easy
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
